# install/ubuntu/systemd/trashyneighbors.service
[Unit]
Description=TrashyNeighbors main site
After=network.target mariadb.service

[Service]
Type=simple
User=trashyneighbors
Group=trashyneighbors
WorkingDirectory=/opt/trashyneighbors
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/trashyneighbors/venv/bin/gunicorn -w 3 -b 127.0.0.1:8000 wsgi:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
